# ansible/roles/bento-icdc/templates/app.timer.j2
[Unit]
Description=Run foo weekly and on boot

[Timer]
OnStartupSec=120

[Install]
WantedBy=timers.target